# dotfiles/update-and-reboot_arch.service
[Unit]
Description=Update packages and reboot

[Service]
Type=exec
ExecStart=pacman -Syu --noconfirm
ExecStartPost=/usr/bin/systemctl reboot
